# crates/server-deck/scripts/network-deck-server.service
# systemd unit for the Deck-side usbip control plane.
#
# Install:
#   sudo cp network-deck-server.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now network-deck-server.service
#
# Depends on usbipd.service (provided by the `usbip` arch package) being up
# so `usbip bind` has a daemon to register the controller with.
#
# First-time setup: pair with the Windows PC by running
#   sudo systemctl stop network-deck-server.service
#   sudo -u deck /usr/local/bin/server-deck pair --state-dir /var/lib/network-deck
# while `client-win` is in pair mode on the PC, then re-enable the unit.

[Unit]
Description=Network Deck usbip control plane (Deck -> Windows)
Documentation=https://github.com/ForbesGRyan/network-steam-deck-controller
After=network-online.target usbipd.service
Wants=network-online.target
Requires=usbipd.service

[Service]
Type=simple
User=root
Group=root

StateDirectory=network-deck

ExecStart=/usr/local/bin/server-deck --state-dir /var/lib/network-deck

Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
