# crates/server-deck/scripts/network-deck-server.service
# systemd unit for the Deck-side input bridge.
#
# Install:
#   sudo cp network-deck-server.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now network-deck-server.service
#
# First-time setup: pair with the Windows PC by running
#   sudo systemctl stop network-deck-server.service
#   /usr/local/bin/server-deck pair /dev/hidraw3 --state-dir /var/lib/network-deck
# while `client-win pair` is running on the PC, then re-enable the unit.

[Unit]
Description=Network Deck input bridge (Steam Deck -> Windows)
Documentation=https://github.com/ForbesGRyan/network-steam-deck-controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deck
Group=deck

Environment=DECK_HIDRAW=/dev/hidraw3
StateDirectory=network-deck

ExecStart=/usr/local/bin/server-deck ${DECK_HIDRAW} --state-dir /var/lib/network-deck

Restart=on-failure
RestartSec=5

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
